# teltonika-django/teltonika-django.service
[Unit]
Description=Teltonika Django GPS Tracking API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=teltonika
Group=teltonika
WorkingDirectory=/opt/teltonika/django
ExecStart=/opt/teltonika/venv/bin/gunicorn --bind 127.0.0.1:8000 --workers 4 --timeout 120 --access-logfile /var/log/teltonika/gunicorn-access.log --error-logfile /var/log/teltonika/gunicorn-error.log --log-level info teltonika_gps.wsgi:application

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Performance settings
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
