# teltonika-django/teltonika-django.service
[Unit]
Description=Teltonika Django GPS Tracking API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=teltonika
Group=teltonika

WorkingDirectory=/opt/teltonika/django
ExecStart=/opt/teltonika/django/venv/bin/gunicorn --bind 0.0.0.0:8000 --workers 2 --timeout 120 --access-logfile /var/log/teltonika/gunicorn-access.log --error-logfile /var/log/teltonika/gunicorn-error.log --log-level info teltonika_gps.wsgi:application

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
